# config/systemd/cloudflared.service
[Unit]
Description=Cloudflare Tunnel
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
# Token 和路径由安装脚本动态替换
ExecStart=__INSTALL_DIR__/bin/cloudflared tunnel --no-autoupdate run --token YOUR_TUNNEL_TOKEN

Restart=on-failure
RestartSec=10
TimeoutStopSec=30

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cloudflared

[Install]
WantedBy=multi-user.target
